# monitor-bootc-version.service
[Unit]
Description=Monitor Fedora bootc container version
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tumeya/git/f-bootc-ws
ExecStart=/usr/bin/python3 /home/tumeya/git/f-bootc-ws/monitor_bootc_version.py

# Run as the user who owns the repository
User=tumeya
Group=tumeya

# Environment variables for git operations
Environment="HOME=/home/tumeya"
Environment="GIT_SSH_COMMAND=ssh -i /home/tumeya/.ssh/id_rsa -o IdentitiesOnly=yes -F /home/tumeya/.ssh/config -o StrictHostKeyChecking=accept-new"
Environment="GIT_AUTHOR_NAME=tumeya"
Environment="GIT_AUTHOR_EMAIL=umeyatakuma@tutanota.com"
Environment="GIT_COMMITTER_NAME=tumeya"
Environment="GIT_COMMITTER_EMAIL=umeyatakuma@tutanota.com"

# Security hardening (optional)
PrivateTmp=yes
#NoNewPrivileges=yes

[Install]
WantedBy=default.target
